# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:edge-tts]
directory=/app/tts
command=/bin/bash -c "export API_KEY='1' && python3 app/server.py"
autostart=true
autorestart=true
environment=PORT=5050
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# DS2API 服务（Go 二进制，监听 8000，自动继承 DS2API_ 环境变量）
[program:deepseek]
directory=/app/deepseek
command=/app/start-ds2api.sh
autostart=true
autorestart=true
environment=PORT=8000,TZ=Asia/Shanghai
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# QwenChat2Api 服务（Node，监听 8001）
[program:qw]
directory=/app/qw
command=npm start
autostart=true
autorestart=true
environment=PORT=8001
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Qwen2API 服务（Node，监听 3000，使用过滤脚本隔离 DS2API 变量）
[program:qwen]
directory=/app/qwen
command=/app/start-qwen.sh
autostart=true
autorestart=true
environment=PORT=3000,SERVICE_PORT=3000,LISTEN_ADDRESS=0.0.0.0,CACHE_MODE=file,PM2_INSTANCES=1,NODE_OPTIONS="--max-old-space-size=160"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
